# rootfs/etc/s6-overlay/s6-rc.d/config-start/up
# This is present only to provide an event. It does nothing else